# supervisor/hbase_big_data.ini
[program:hbase_big_data]
command=/home/anaconda3/bin/python3  /home/hadoop/hbase_big_data/app.py
stopasgroup=true
killasgroup=true
directory=/home/hadoop/hbase_big_data/
autostart=true
autorestart=true
user=root
redirect_stderr=true
stdout_logfile=/home/log_sellbot/hbase_big_data/supervisord-10002.log