# supervisor/queryserver.ini
[program:queryserver]
command=/home/hadoop/phoenix/bin/queryserver.py start
stopasgroup=true
killasgroup=true
directory=/home/hadoop/hbase_big_data/
autostart=true
autorestart=true
user=root
redirect_stderr=true
stdout_logfile=/home/log_sellbot/hbase_big_data/supervisord-queryserver10001.log